# debian/status-area-tor.postinst
#!/bin/sh

gtk-update-icon-cache -f /usr/share/icons/hicolor
rc-update del tor default || true

mkdir -p /var/lib/tor-applet
chown debian-tor:debian-tor /var/lib/tor-applet

#DEBHELPER#
