# debian/status-area-tor.postrm
#!/bin/sh

update-rc.d tor enable || true

#DEBHELPER#
